# .config/systemd/user/coldsync-metrics.service
[Unit]
Description=Emit coldsync metrics

[Service]
Type=oneshot
ExecStart=%h/DuRiWorkspace/scripts/evolution/coldsync_emit_metrics.sh

