# .config/systemd/user/duri-synth.timer
[Unit]
Description=Hourly synthetic probe
[Timer]
OnCalendar=hourly
[Install]
WantedBy=default.target
