# .config/systemd/user/kimshin-dr-restore.timer
[Unit]
Description=Run DR restore rehearsal monthly
[Timer]
OnCalendar=monthly
RandomizedDelaySec=15m
Persistent=true
[Install]
WantedBy=default.target
