# .config/systemd/user/kimshin-dr-restore.service
[Unit]
Description=Monthly DR restore rehearsal (byte-identical)
[Service]
Type=oneshot
WorkingDirectory=%h/DuRiWorkspace
Environment=UMASK=077
ExecStart=/usr/bin/bash -lc 'backup/verify_and_restore.sh && echo "[DR] restore OK"'
