# .config/systemd/user/kimshin-backup.service
[Unit]
Description=Secure backup for KIMSHIN charter

[Service]
Type=oneshot
WorkingDirectory=%h/DuRiWorkspace
Environment=AGE_PASSPHRASE=***
ExecStart=%h/DuRiWorkspace/backup/secure_backup.sh
