# systemd/duri-hygiene.timer
[Unit]
Description=Run data hygiene daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
