# systemd/duri-metrics-day66.timer
[Unit]
Description=Run Day66 metrics every hour (at minute 10)

[Timer]
OnCalendar=*:10
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target


